# test/files/sleep@.service
[Unit]
Description=Sleeps forever

[Service]
Restart=always
ExecStart=sleep infinity

[Install]
WantedBy=multi-user.target
